# etc/supervisor/conf.d/davmail.conf
[program:DavMail]
command=/opt/davmail/davmail.sh /etc/davmail.properties
user=root
autostart=true
autorestart=true
startsecs=120
stopasgroup=true
ikillasgroup=true
startretries=1
